# services/slitter_restart.timer
[Unit]
Description=auto slitter servo restart timer

[Timer]
OnActiveSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target