# services/slitter_restart.service
[Unit]
Description=restart auto slitter servo control

[Service]
ExecStart=/bin/systemctl restart scon